# setup-server.service
[Unit]
After=network.target

[Service]
ExecStart=/home/aidan/home-server/setup.sh

[Install]
WantedBy=multi-user.target